# server-backup.cron
# Server Backup
SHELL=/bin/bash
PATH=/sbin:/bin:/usr/sbin:/usr/bin
MAILTO=youremail@yourdomain.com
30 02 * * * root /usr/local/bin/server-backup.sh

